# internal/install/templates/guild-daemon.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  guild daemon launchd agent. Generated by 'guild daemon install';
  re-run that command after moving the guild binary, or remove with
  'guild daemon uninstall'.

  No StandardOutPath/StandardErrorPath keys on purpose: the daemon
  writes its own log to ~/.guild/daemon.log, so launchd does not need
  to capture stdout/stderr.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>{{.Label}}</string>
	<key>ProgramArguments</key>
	<array>
		<string>{{.BinPath}}</string>
		<string>daemon</string>
		<string>run</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
